# board/pluto/overlay/etc/profile.d/plutosdr_prompt.sh
export PS1="\[\e[0:34m\]pluto\[\e[m\]\[\e[1;34m\]sdr\[\e[m\]\[\e[0;34m\]>\[\e[m\] "
